# php/7.2-fpm/conf.d/supervisor/nginx.conf
[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=nginx
autostart=true
autorestart=true
user=nginx
numprocs=2
redirect_stderr=true